# addons/ha-aux/rootfs/etc/services.d/ha-aux/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting HA aux Audio Renderer..."

exec python3 -m ha_aux_addon
